# systemd/assassind.service
[Unit]
Name=Assassind daemon

[Service]
ExecStart=/usr/bin/assassind

[Install]
WantedBy=multi-user.target

