# systemd/oswl_keystone_user_collectord.service
[Unit]
Name=oswl keystone user collector daemon

[Service]
ExecStart=/usr/bin/oswl_collectord keystone_user

[Install]
WantedBy=multi-user.target

