# systemd/oswl_flavor_collectord.service
[Unit]
Name=oswl flavor collector daemon

[Service]
ExecStart=/usr/bin/oswl_collectord flavor

[Install]
WantedBy=multi-user.target

