# systemd/oswl_volume_collectord.service
[Unit]
Name=oswl volume collector daemon

[Service]
ExecStart=/usr/bin/oswl_collectord volume

[Install]
WantedBy=multi-user.target

