# systemd/oswl_tenant_collectord.service
[Unit]
Name=oswl tenant collector  daemon

[Service]
ExecStart=/usr/bin/oswl_collectord tenant

[Install]
WantedBy=multi-user.target

