# systemd/receiverd.service
[Unit]
Name=Receiverd daemon

[Service]
ExecStart=/usr/bin/receiverd

[Install]
WantedBy=multi-user.target

