# systemd/statsenderd.service
[Unit]
Name=Statsenderd daemon

[Service]
ExecStart=/usr/bin/statsenderd

[Install]
WantedBy=multi-user.target

